# files/rootfs/lenny_working_fs/var/lib/dpkg/info/module-init-tools.conffiles
/etc/modprobe.d/aliases
/etc/init.d/module-init-tools
