# lenny_sandbox/var/lib/dpkg/info/libc6.conffiles
/etc/bindresvport.blacklist
/etc/gai.conf
/etc/ld.so.conf.d/x86_64-linux-gnu.conf
/etc/ld.so.conf.d/libc.conf
/etc/init.d/glibc.sh
