# lenny_sandbox/var/lib/dpkg/info/module-init-tools.conffiles
/etc/init.d/module-init-tools
/etc/modprobe.d/aliases
/etc/modprobe.d/arch/x86_64
