# lenny_sandbox/var/lib/dpkg/info/rsyslog.conffiles
/etc/init.d/rsyslog
/etc/default/rsyslog
/etc/logrotate.d/rsyslog
/etc/rsyslog.conf
