# lenny_sandbox/var/lib/dpkg/info/uml-utilities.conffiles
/etc/init.d/uml-utilities
/etc/default/uml-utilities
/etc/network/if-pre-up.d/uml-utilities
/etc/network/if-up.d/uml-utilities
