# lenny_sandbox/var/lib/dpkg/info/uml-utilities.list
/.
/etc
/etc/init.d
/etc/init.d/uml-utilities
/etc/default
/etc/default/uml-utilities
/etc/network
/etc/network/if-pre-up.d
/etc/network/if-pre-up.d/uml-utilities
/etc/network/if-up.d
/etc/network/if-up.d/uml-utilities
/var
/var/run
/var/run/uml-utilities
/usr
/usr/sbin
/usr/sbin/tunctl
/usr/sbin/jail_uml
/usr/lib
/usr/lib/uml
/usr/lib/uml/uml_net
/usr/lib/uml/port-helper
/usr/bin
/usr/bin/uml_watchdog
/usr/bin/humfsify
/usr/bin/uml_mconsole
/usr/bin/uml_switch
/usr/bin/uml_moo
/usr/bin/uml_mkcow
/usr/bin/uml_mount
/usr/bin/jailtest
/usr/share
/usr/share/doc
/usr/share/doc/uml-utilities
/usr/share/doc/uml-utilities/copyright
/usr/share/doc/uml-utilities/README.Debian
/usr/share/doc/uml-utilities/examples
/usr/share/doc/uml-utilities/examples/interfaces.example
/usr/share/doc/uml-utilities/examples/tty_log.pl
/usr/share/doc/uml-utilities/changelog.Debian.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/uml_switch.1.gz
/usr/share/man/man1/uml_mount.1.gz
/usr/share/man/man1/jail_uml.1.gz
/usr/share/man/man1/humfsify.1.gz
/usr/share/man/man1/uml_moo.1.gz
/usr/share/man/man1/uml_mkcow.1.gz
/usr/share/man/man1/uml_mconsole.1.gz
/usr/share/man/man8
/usr/share/man/man8/tunctl.8.gz
