# lenny_sandbox/var/lib/dpkg/info/util-vserver.list
/.
/etc
/etc/init.d
/etc/init.d/util-vserver
/etc/default
/etc/default/util-vserver
/etc/vservers
/etc/vservers/.defaults
/etc/vservers/.defaults/files
/etc/vservers/.defaults/apps
/etc/vservers/.defaults/apps/vunify
/etc/vservers/.defaults/apps/vunify/hash
/etc/vservers/.distributions
/etc/vservers/.distributions/fc4
/etc/vservers/.distributions/fc4/apt
/etc/vservers/.distributions/fc4/apt/sources.list
/etc/vservers/.distributions/fc5
/etc/vservers/.distributions/fc5/apt
/etc/vservers/.distributions/fc5/apt/sources.list
/etc/vservers/.distributions/f7
/etc/vservers/.distributions/f7/apt
/etc/vservers/.distributions/f7/apt/sources.list
/etc/vservers/.distributions/suse91
/etc/vservers/.distributions/suse91/apt
/etc/vservers/.distributions/suse91/apt/sources.list
/etc/vservers/.distributions/.common
/etc/vservers/.distributions/.common/pubkeys
/etc/vservers/.distributions/fc2
/etc/vservers/.distributions/fc2/apt
/etc/vservers/.distributions/fc2/apt/sources.list
/etc/vservers/.distributions/fc6
/etc/vservers/.distributions/fc6/apt
/etc/vservers/.distributions/fc6/apt/sources.list
/etc/vservers/.distributions/fc1
/etc/vservers/.distributions/fc1/apt
/etc/vservers/.distributions/fc1/apt/sources.list
/etc/vservers/.distributions/rh9
/etc/vservers/.distributions/rh9/apt
/etc/vservers/.distributions/rh9/apt/sources.list
/etc/vservers/.distributions/fc3
/etc/vservers/.distributions/fc3/apt
/etc/vservers/.distributions/fc3/apt/sources.list
/etc/bash_completion.d
/etc/bash_completion.d/vserver
/var
/var/lock
/var/lock/subsys
/var/cache
/var/cache/vservers
/var/lib
/var/lib/vservers
/var/lib/vservers/.pkg
/var/run
/var/run/vservers.rev
/var/run/vservers
/var/run/vshelper
/sbin
/usr
/usr/sbin
/usr/sbin/vuname
/usr/sbin/vattribute
/usr/sbin/ncontext
/usr/sbin/vrsetup
/usr/sbin/vesync
/usr/sbin/vpstree
/usr/sbin/vdevmap
/usr/sbin/vps
/usr/sbin/vrpm
/usr/sbin/vsched
/usr/sbin/vtag
/usr/sbin/chbind
/usr/sbin/vkill
/usr/sbin/vapt-get
/usr/sbin/chxid
/usr/sbin/vsomething
/usr/sbin/vdispatch-conf
/usr/sbin/chcontext
/usr/sbin/vnamespace
/usr/sbin/exec-cd
/usr/sbin/vserver-info
/usr/sbin/vserver
/usr/sbin/vupdateworld
/usr/sbin/naddress
/usr/sbin/reducecap
/usr/sbin/vwait
/usr/sbin/vdu
/usr/sbin/setattr
/usr/sbin/vmount
/usr/sbin/vserver-stat
/usr/sbin/lsxid
/usr/sbin/vdlimit
/usr/sbin/vmemctrl
/usr/sbin/nattribute
/usr/sbin/vyum
/usr/sbin/vtop
/usr/sbin/vcontext
/usr/sbin/vspace
/usr/sbin/vlimit
/usr/sbin/showattr
/usr/sbin/vemerge
/usr/lib
/usr/lib/libvserver.a
/usr/lib/util-vserver
/usr/lib/util-vserver/vserver-wrapper
/usr/lib/util-vserver/secure-mount
/usr/lib/util-vserver/vserver-build.apt-rpm
/usr/lib/util-vserver/vrpm-worker
/usr/lib/util-vserver/defaults
/usr/lib/util-vserver/defaults/fstab
/usr/lib/util-vserver/defaults/vprocunhide-files
/usr/lib/util-vserver/defaults/context.start
/usr/lib/util-vserver/defaults/mtab
/usr/lib/util-vserver/defaults/environment
/usr/lib/util-vserver/defaults/vunify-exclude
/usr/lib/util-vserver/defaults/h2ext.desc
/usr/lib/util-vserver/defaults/debootstrap.uri
/usr/lib/util-vserver/defaults/debootstrap.mirror
/usr/lib/util-vserver/keep-ctx-alive
/usr/lib/util-vserver/vserver.stop
/usr/lib/util-vserver/vunify
/usr/lib/util-vserver/lockfile
/usr/lib/util-vserver/vserver-build.functions.rpm
/usr/lib/util-vserver/pkgmgmt
/usr/lib/util-vserver/sigexec
/usr/lib/util-vserver/vclone
/usr/lib/util-vserver/vserver-build.yum
/usr/lib/util-vserver/chain-echo
/usr/lib/util-vserver/mask2prefix
/usr/lib/util-vserver/vserver-build.functions.pkgmgmt
/usr/lib/util-vserver/vapt-get-worker
/usr/lib/util-vserver/vserver-build.rpm
/usr/lib/util-vserver/vsysctl
/usr/lib/util-vserver/vserver-build.skeleton
/usr/lib/util-vserver/vyum-worker
/usr/lib/util-vserver/readlink
/usr/lib/util-vserver/vserver-build.functions
/usr/lib/util-vserver/vserver.suexec
/usr/lib/util-vserver/tunctl
/usr/lib/util-vserver/vshelper
/usr/lib/util-vserver/vservers.grabinfo.sh
/usr/lib/util-vserver/vserver-build.rsync
/usr/lib/util-vserver/vserver-init.functions
/usr/lib/util-vserver/vshelper-sync
/usr/lib/util-vserver/save_ctxinfo
/usr/lib/util-vserver/capchroot
/usr/lib/util-vserver/check-unixfile
/usr/lib/util-vserver/vserver-build.fai
/usr/lib/util-vserver/vserver-build
/usr/lib/util-vserver/h2ext
/usr/lib/util-vserver/h2ext-worker
/usr/lib/util-vserver/vserver-setup.functions
/usr/lib/util-vserver/vpkg
/usr/lib/util-vserver/vserver-build.functions.yum
/usr/lib/util-vserver/vprocunhide
/usr/lib/util-vserver/vserver.functions
/usr/lib/util-vserver/rpm-fake.so
/usr/lib/util-vserver/functions
/usr/lib/util-vserver/vcopy
/usr/lib/util-vserver/util-vserver-vars
/usr/lib/util-vserver/vserver.start
/usr/lib/util-vserver/exec-remount
/usr/lib/util-vserver/vserver-build.functions.apt
/usr/lib/util-vserver/filetime
/usr/lib/util-vserver/vserver.delete
/usr/lib/util-vserver/FEATURES.txt
/usr/lib/util-vserver/rpm-fake-resolver
/usr/lib/util-vserver/vserver-build.clone
/usr/lib/util-vserver/vserver-build.template
/usr/lib/util-vserver/start-vservers
/usr/lib/util-vserver/vhashify
/usr/lib/util-vserver/vsysvwrapper
/usr/lib/util-vserver/fakerunlevel
/usr/lib/util-vserver/exec-ulimit
/usr/lib/util-vserver/vserver-build.debootstrap
/usr/lib/util-vserver/distributions
/usr/lib/util-vserver/distributions/fc4
/usr/lib/util-vserver/distributions/fc4/yum.repos.d
/usr/lib/util-vserver/distributions/fc4/yum.repos.d/fedora-devel.repo
/usr/lib/util-vserver/distributions/fc4/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/fc4/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/fc4/yum.repos.d/fedora.repo
/usr/lib/util-vserver/distributions/fc4/rpmlist.d
/usr/lib/util-vserver/distributions/fc4/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/fc4/apt
/usr/lib/util-vserver/distributions/fc4/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc4/pkgs
/usr/lib/util-vserver/distributions/fc4/pkgs/02
/usr/lib/util-vserver/distributions/fc4/pkgs/01
/usr/lib/util-vserver/distributions/fc4/pkgs/03
/usr/lib/util-vserver/distributions/fc4/yum
/usr/lib/util-vserver/distributions/fc4/yum/yum.conf
/usr/lib/util-vserver/distributions/fc5
/usr/lib/util-vserver/distributions/fc5/yum.repos.d
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-core.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-legacy.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-extras.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-extras-development.repo
/usr/lib/util-vserver/distributions/fc5/yum.repos.d/fedora-development.repo
/usr/lib/util-vserver/distributions/fc5/rpmlist.d
/usr/lib/util-vserver/distributions/fc5/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/fc5/apt
/usr/lib/util-vserver/distributions/fc5/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc5/pkgs
/usr/lib/util-vserver/distributions/fc5/pkgs/02
/usr/lib/util-vserver/distributions/fc5/pkgs/01
/usr/lib/util-vserver/distributions/fc5/pkgs/03
/usr/lib/util-vserver/distributions/fc5/yum
/usr/lib/util-vserver/distributions/fc5/yum/yum.conf
/usr/lib/util-vserver/distributions/template
/usr/lib/util-vserver/distributions/template/initpost
/usr/lib/util-vserver/distributions/template/initpre
/usr/lib/util-vserver/distributions/centos5
/usr/lib/util-vserver/distributions/centos5/yum.repos.d
/usr/lib/util-vserver/distributions/centos5/yum.repos.d/CentOS-Base.repo
/usr/lib/util-vserver/distributions/centos5/pkgs
/usr/lib/util-vserver/distributions/centos5/pkgs/02
/usr/lib/util-vserver/distributions/centos5/pkgs/01
/usr/lib/util-vserver/distributions/centos5/pkgs/03
/usr/lib/util-vserver/distributions/centos5/yum
/usr/lib/util-vserver/distributions/centos5/yum/yum.conf
/usr/lib/util-vserver/distributions/defaults
/usr/lib/util-vserver/distributions/defaults/rpm
/usr/lib/util-vserver/distributions/defaults/rpm/macros
/usr/lib/util-vserver/distributions/defaults/apt
/usr/lib/util-vserver/distributions/defaults/apt/apt.conf
/usr/lib/util-vserver/distributions/defaults/apt.conf
/usr/lib/util-vserver/distributions/defaults/devs
/usr/lib/util-vserver/distributions/gutsy
/usr/lib/util-vserver/distributions/gutsy/debootstrap.script
/usr/lib/util-vserver/distributions/warty
/usr/lib/util-vserver/distributions/warty/debootstrap.script
/usr/lib/util-vserver/distributions/f7
/usr/lib/util-vserver/distributions/f7/yum.repos.d
/usr/lib/util-vserver/distributions/f7/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/f7/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/f7/yum.repos.d/fedora.repo
/usr/lib/util-vserver/distributions/f7/yum.repos.d/fedora-development.repo
/usr/lib/util-vserver/distributions/f7/rpmlist.d
/usr/lib/util-vserver/distributions/f7/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/f7/apt
/usr/lib/util-vserver/distributions/f7/apt/rpmpriorities
/usr/lib/util-vserver/distributions/f7/pkgs
/usr/lib/util-vserver/distributions/f7/pkgs/02
/usr/lib/util-vserver/distributions/f7/pkgs/01
/usr/lib/util-vserver/distributions/f7/pkgs/03
/usr/lib/util-vserver/distributions/f7/yum
/usr/lib/util-vserver/distributions/f7/yum/yum.conf
/usr/lib/util-vserver/distributions/breezy
/usr/lib/util-vserver/distributions/breezy/debootstrap.script
/usr/lib/util-vserver/distributions/suse91
/usr/lib/util-vserver/distributions/suse91/rpm
/usr/lib/util-vserver/distributions/suse91/rpm/macros
/usr/lib/util-vserver/distributions/suse91/apt
/usr/lib/util-vserver/distributions/suse91/apt/rpmpriorities
/usr/lib/util-vserver/distributions/suse91/pkgs
/usr/lib/util-vserver/distributions/suse91/pkgs/02
/usr/lib/util-vserver/distributions/suse91/pkgs/01
/usr/lib/util-vserver/distributions/hoary
/usr/lib/util-vserver/distributions/hoary/debootstrap.script
/usr/lib/util-vserver/distributions/edgy
/usr/lib/util-vserver/distributions/edgy/debootstrap.script
/usr/lib/util-vserver/distributions/fc2
/usr/lib/util-vserver/distributions/fc2/apt
/usr/lib/util-vserver/distributions/fc2/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc2/pkgs
/usr/lib/util-vserver/distributions/fc2/pkgs/02
/usr/lib/util-vserver/distributions/fc2/pkgs/01
/usr/lib/util-vserver/distributions/redhat
/usr/lib/util-vserver/distributions/redhat/initctl
/usr/lib/util-vserver/distributions/redhat/initpost
/usr/lib/util-vserver/distributions/redhat/rc.sysinit
/usr/lib/util-vserver/distributions/redhat/initpre
/usr/lib/util-vserver/distributions/gentoo
/usr/lib/util-vserver/distributions/gentoo/shutdown.sh
/usr/lib/util-vserver/distributions/gentoo/init-vserver.sh
/usr/lib/util-vserver/distributions/gentoo/initpost
/usr/lib/util-vserver/distributions/gentoo/net.vserver
/usr/lib/util-vserver/distributions/gentoo/reboot.sh
/usr/lib/util-vserver/distributions/gentoo/initpre
/usr/lib/util-vserver/distributions/alpine
/usr/lib/util-vserver/distributions/alpine/initpost
/usr/lib/util-vserver/distributions/alpine/initpre
/usr/lib/util-vserver/distributions/f8
/usr/lib/util-vserver/distributions/f8/yum.repos.d
/usr/lib/util-vserver/distributions/f8/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/f8/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/f8/yum.repos.d/fedora.repo
/usr/lib/util-vserver/distributions/f8/yum.repos.d/fedora-development.repo
/usr/lib/util-vserver/distributions/f8/rpmlist.d
/usr/lib/util-vserver/distributions/f8/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/f8/apt
/usr/lib/util-vserver/distributions/f8/apt/rpmpriorities
/usr/lib/util-vserver/distributions/f8/pkgs
/usr/lib/util-vserver/distributions/f8/pkgs/02
/usr/lib/util-vserver/distributions/f8/pkgs/01
/usr/lib/util-vserver/distributions/f8/pkgs/03
/usr/lib/util-vserver/distributions/f8/yum
/usr/lib/util-vserver/distributions/f8/yum/yum.conf
/usr/lib/util-vserver/distributions/feisty
/usr/lib/util-vserver/distributions/feisty/debootstrap.script
/usr/lib/util-vserver/distributions/centos4
/usr/lib/util-vserver/distributions/centos4/yum.repos.d
/usr/lib/util-vserver/distributions/centos4/yum.repos.d/CentOS-Base.repo
/usr/lib/util-vserver/distributions/centos4/pkgs
/usr/lib/util-vserver/distributions/centos4/pkgs/02
/usr/lib/util-vserver/distributions/centos4/pkgs/01
/usr/lib/util-vserver/distributions/centos4/pkgs/03
/usr/lib/util-vserver/distributions/centos4/yum
/usr/lib/util-vserver/distributions/centos4/yum/yum.conf
/usr/lib/util-vserver/distributions/fc6
/usr/lib/util-vserver/distributions/fc6/yum.repos.d
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-core.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-legacy.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-extras.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-extras-development.repo
/usr/lib/util-vserver/distributions/fc6/yum.repos.d/fedora-development.repo
/usr/lib/util-vserver/distributions/fc6/rpmlist.d
/usr/lib/util-vserver/distributions/fc6/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/fc6/apt
/usr/lib/util-vserver/distributions/fc6/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc6/pkgs
/usr/lib/util-vserver/distributions/fc6/pkgs/02
/usr/lib/util-vserver/distributions/fc6/pkgs/01
/usr/lib/util-vserver/distributions/fc6/pkgs/03
/usr/lib/util-vserver/distributions/fc6/yum
/usr/lib/util-vserver/distributions/fc6/yum/yum.conf
/usr/lib/util-vserver/distributions/debian
/usr/lib/util-vserver/distributions/debian/debootstrap.script
/usr/lib/util-vserver/distributions/debian/initpost
/usr/lib/util-vserver/distributions/fc1
/usr/lib/util-vserver/distributions/fc1/apt
/usr/lib/util-vserver/distributions/fc1/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc1/pkgs
/usr/lib/util-vserver/distributions/fc1/pkgs/02
/usr/lib/util-vserver/distributions/fc1/pkgs/01
/usr/lib/util-vserver/distributions/rh9
/usr/lib/util-vserver/distributions/rh9/apt
/usr/lib/util-vserver/distributions/rh9/apt/rpmpriorities
/usr/lib/util-vserver/distributions/rh9/pkgs
/usr/lib/util-vserver/distributions/rh9/pkgs/02
/usr/lib/util-vserver/distributions/rh9/pkgs/01
/usr/lib/util-vserver/distributions/fc3
/usr/lib/util-vserver/distributions/fc3/yum.repos.d
/usr/lib/util-vserver/distributions/fc3/yum.repos.d/fedora-devel.repo
/usr/lib/util-vserver/distributions/fc3/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/fc3/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/fc3/yum.repos.d/fedora.repo
/usr/lib/util-vserver/distributions/fc3/rpmlist.d
/usr/lib/util-vserver/distributions/fc3/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/fc3/apt
/usr/lib/util-vserver/distributions/fc3/apt/rpmpriorities
/usr/lib/util-vserver/distributions/fc3/pkgs
/usr/lib/util-vserver/distributions/fc3/pkgs/02
/usr/lib/util-vserver/distributions/fc3/pkgs/01
/usr/lib/util-vserver/distributions/fc3/pkgs/03
/usr/lib/util-vserver/distributions/fc3/yum
/usr/lib/util-vserver/distributions/fc3/yum/yum.conf
/usr/lib/util-vserver/distributions/dapper
/usr/lib/util-vserver/distributions/dapper/debootstrap.script
/usr/lib/util-vserver/distributions/f9
/usr/lib/util-vserver/distributions/f9/yum.repos.d
/usr/lib/util-vserver/distributions/f9/yum.repos.d/fedora-updates-testing.repo
/usr/lib/util-vserver/distributions/f9/yum.repos.d/fedora-updates.repo
/usr/lib/util-vserver/distributions/f9/yum.repos.d/fedora.repo
/usr/lib/util-vserver/distributions/f9/yum.repos.d/fedora-development.repo
/usr/lib/util-vserver/distributions/f9/rpmlist.d
/usr/lib/util-vserver/distributions/f9/rpmlist.d/00.lst
/usr/lib/util-vserver/distributions/f9/pkgs
/usr/lib/util-vserver/distributions/f9/pkgs/02
/usr/lib/util-vserver/distributions/f9/pkgs/01
/usr/lib/util-vserver/distributions/f9/pkgs/03
/usr/lib/util-vserver/distributions/f9/yum
/usr/lib/util-vserver/distributions/f9/yum/yum.conf
/usr/lib/util-vserver/chroot-sh
/usr/lib/util-vserver/chbind-compat
/usr/lib/util-vserver/vrpm-preload
/usr/lib/util-vserver/chcontext-compat
/usr/lib/pkgconfig
/usr/lib/pkgconfig/util-vserver.pc
/usr/lib/libvserver.so.0.0.0
/usr/include
/usr/include/vserver.h
/usr/share
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/util-vserver
/usr/share/doc
/usr/share/doc/util-vserver
/usr/share/doc/util-vserver/copyright
/usr/share/doc/util-vserver/NEWS.gz
/usr/share/doc/util-vserver/README.Debian
/usr/share/doc/util-vserver/README.gz
/usr/share/doc/util-vserver/examples
/usr/share/doc/util-vserver/examples/README.TestScripts
/usr/share/doc/util-vserver/examples/testme.sh-0.17.gz
/usr/share/doc/util-vserver/examples/vschedcalc
/usr/share/doc/util-vserver/examples/testfs.sh-0.18.gz
/usr/share/doc/util-vserver/changelog.Debian.gz
/usr/share/doc/util-vserver/changelog.gz
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/vserver-stat.8.gz
/usr/share/man/man8/chcontext.8.gz
/usr/share/man/man8/vserver.8.gz
/usr/share/man/man8/chbind.8.gz
/usr/share/man/man8/vtop.8.gz
/usr/share/man/man8/vpstree.8.gz
/usr/share/man/man8/vserver-build.8.gz
/usr/share/man/man8/reducecap.8.gz
/etc/vservers/.defaults/cachebase
/sbin/vshelper
/usr/lib/libvserver.so.0
/usr/lib/util-vserver/distributions/fc4/initpost
/usr/lib/util-vserver/distributions/fc4/initpre
/usr/lib/util-vserver/distributions/fc5/initpost
/usr/lib/util-vserver/distributions/fc5/initpre
/usr/lib/util-vserver/distributions/centos5/initpost
/usr/lib/util-vserver/distributions/centos5/initpre
/usr/lib/util-vserver/distributions/f7/initpost
/usr/lib/util-vserver/distributions/f7/initpre
/usr/lib/util-vserver/distributions/sid
/usr/lib/util-vserver/distributions/etch
/usr/lib/util-vserver/distributions/fc2/initpost
/usr/lib/util-vserver/distributions/fc2/initpre
/usr/lib/util-vserver/distributions/f8/initpost
/usr/lib/util-vserver/distributions/f8/initpre
/usr/lib/util-vserver/distributions/centos4/initpost
/usr/lib/util-vserver/distributions/centos4/initpre
/usr/lib/util-vserver/distributions/fc6/initpost
/usr/lib/util-vserver/distributions/fc6/initpre
/usr/lib/util-vserver/distributions/fc1/initpost
/usr/lib/util-vserver/distributions/fc1/initpre
/usr/lib/util-vserver/distributions/rh9/initpost
/usr/lib/util-vserver/distributions/rh9/initpre
/usr/lib/util-vserver/distributions/fc3/initpost
/usr/lib/util-vserver/distributions/fc3/initpre
/usr/lib/util-vserver/distributions/f9/initpost
/usr/lib/util-vserver/distributions/f9/initpre
/usr/lib/util-vserver/distributions/lenny
/usr/lib/libvserver.so
